# server/extras/flink.service
[Unit]
Description=Flink internal prototyping server
Documentation=https://github.com/csweichel/flink
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/flink
ExecStart=%h/.local/bin/flink-server --config %h/.config/flink/flink.yaml
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
